# services/jetson-full-fan.service
[Unit]
Description=Set fan at max

[Service]
Type=oneshot
ExecStart=/bin/sh -c "/usr/bin/jetson_clocks --fan"

[Install]
WantedBy=multi-user.target
