# services/limit-system-cpus.service
[Unit]
Description=Limit CPUs that system processes can use

[Service]
Type=oneshot
ExecStart=/bin/bash /startup_scripts/limit-system-cpus.sh
